# pyrouted.socket.in
[Socket]
ListenStream=/var/run/pyrouted.sock
SocketUser=root
SocketGroup=wheel
SocketMode=660
Accept=true
