# pyrouted.service.in
[Service]
Environment=@ENVIRON@
ExecStart=@PYROUTED@
User=@USER@
StandardOutput=syslog
StandardError=syslog
