# pyrouted.service.in
[Service]
Environment=%ENVIRON%
ExecStart=%PYROUTED% %SOCKET%
User=%USER%
StandardOutput=syslog
StandardError=syslog
